# templates/etc/systemd/system/apt-daily.timer.d/override.conf.j2
# {{ ansible_managed }}
# ansibleguy.linux_bootstrap

[Timer]
OnCalendar=
OnCalendar={{ BOOTSTRAP.auto_update.update_day }} {{ BOOTSTRAP.auto_update.time_update }}
