# templates/etc/systemd/system/apt-daily.timer.d/override.conf.j2
[Timer]
OnCalendar=
OnCalendar={{ auto_update_config.time }}
